# deploy/vps/alpaca-paper-bot-main.timer
[Unit]
Description=Start Alpaca Paper Bot main session guard at 16:15 Europe/Istanbul

[Timer]
OnCalendar=Mon..Fri 16:15:00
AccuracySec=1s
RandomizedDelaySec=0
Persistent=true
Unit=alpaca-paper-bot-main.service

[Install]
WantedBy=timers.target
